# versions/bookworm.v3/package/etc/supervisor.d/cron.ini
[program:cron]
environment=AWS_ACCOUNT_ID="%(ENV_AWS_ACCOUNT_ID)s",
  AWS_DEFAULT_REGION="%(ENV_AWS_DEFAULT_REGION)s",
  APIGATEWAY_NAME="%(ENV_APIGATEWAY_NAME)s",
  APIGATEWAY_API_KEY_NAME="%(ENV_APIGATEWAY_API_KEY_NAME)s",
  APIGATEWAY_STAGE="%(ENV_APIGATEWAY_STAGE)s"
command=/usr/sbin/cron -f -L 15
autostart=true
autorestart=true