# versions/bookworm.v3/package/etc/supervisor.d/cron.ini
[program:cron]
environment=AWS_ACCOUNT_ID="%(ENV_AWS_ACCOUNT_ID)s"
command=/usr/sbin/cron -f -L 15
autostart=true
autorestart=true