# 0.0.1/etc/supervisor.d/ttyd.ini
[program:ttyd]
command = /usr/bin/ttyd -p 8080 /bin/bash
autostart = true
autorestart = true
stderr_logfile = /var/log/ttyd.err.log
stdout_logfile = /var/log/ttyd.out.log