# 0.0.1/etc/supervisor.d/ttyd.ini
[program:ttyd]

command = /usr/bin/ttyd \
    --socket-owner=nginx:nginx \
    --interface=/run/nginx/ttyd.socket \
    --writable \
    --uid=1000 --gid=1000 /bin/bash
autostart = true
autorestart = true
stderr_logfile = /var/log/ttyd.err.log
stdout_logfile = /var/log/ttyd.out.log