# versions/bookworm.v3/package/etc/cron.d/awsmgr
## Note cannot use */1 -> * is every min
*/5 * * * * cyverse test -e /home/cyverse/envs/flask-env/bin/awsmgr && /home/cyverse/envs/flask-env/bin/awsmgr renew-token 2>&1 | logger -t awsmgr -p cron.err
