# versions/bookworm.v3/package/etc/cron.d/deadman
*/5 * * * * cyverse /usr/local/bin/poke-deadman.sh 2>&1 | logger -t awsmgr -p cron.err