# versions/bookworm.v3/package/etc/cron.d/deadman
*/5 * * * * cyverse /usr/local/bin/poke-deadman.sh | logger -t deadman -p cron.err
