# versions/bookworm.v2/etc/supervisor.d/gunicorn.ini
[program:gunicorn]

command = /home/gunicorn/envs/flask-env/bin/gunicorn -c /etc/gunicorn/awsmgr.conf.py awsmgr.wsgi:app
directory = /usr/src/awsmgr
user = gunicorn
group = gunicorn
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/gunicorn.log