# versions/bookworm.v2/etc/supervisor.d/ttyd.ini
[program:ttyd]
environment=AWS_KMS_KEY="%(ENV_AWS_KMS_KEY)s",HOME=/home/cyverse
command=/usr/local/bin/ttyd -p 7681 --uid=1000 --gid=1000 --writable --check-origin /usr/bin/bash
autostart = true
autorestart = true
stderr_logfile = /var/log/ttyd/error.log
stdout_logfile = /var/log/ttyd/access.log